# scripts/net/etc/cron.d/gyrid-scripts-net
# Call the NTP restart script every 5 minutes (only restart if necessary)
*/5 * * * * root /usr/sbin/mobile-internet-restart &> /dev/null
