# scripts/ntp/etc/cron.d/gyrid-scripts-ntp
# Call the NTP restart script every 15 minutes (only restart if necessary)
3,18,33,48 * * * * root /usr/sbin/ntp-restart &> /dev/null
