# scripts/debian/gyrid-scripts-ssh.postrm
#!/bin/sh

if [ "$1" = purge ]; then
    update-rc.d -f ssh-phone-home remove
fi
