# scripts/debian/gyrid-scripts-core.postinst
#!/bin/sh

PKG=gyrid-scripts-net

update-rc.d gyrid-os defaults 99 01
