# scripts/debian/gyrid-scripts-core.prerm
#!/bin/sh

set -e

if which invoke-rc.d >/dev/null 2>&1; then
    invoke-rc.d gyrid-os stop
else
    /etc/init.d/gyrid-os stop
fi
