# scripts/debian/gyrid-scripts-core.postrm
#!/bin/sh

if [ "$1" = purge ]; then
    update-rc.d -f gyrid-os remove
fi
