# scripts/core/etc/cron.d/gyrid-scripts-core
* * * * * root /usr/sbin/gyrid-led-checker &> /dev/null
* * * * * root /usr/sbin/gyrid-crash-restart &> /dev/null
