# scripts/apt/etc/cron.d/cron-apt
#
# Regular cron jobs for the cron-apt package
#
0 *	* * *	root	test -x /usr/sbin/cron-apt && /usr/sbin/cron-apt /etc/cron-apt/bluemap.conf
#0 4	* * *	root	test -x /usr/sbin/cron-apt && /usr/sbin/cron-apt /etc/cron-apt/debian.conf

