# scripts/ssh/etc/cron.d/gyrid-scripts-ssh
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Call the autoSSH restart script every 15 minutes (only restart if necessary)
4,19,34,49 * * * * root /usr/sbin/autossh-restart &> /dev/null
